# geosync/georchestra-ouvert/crontab/georchestra-ouvert
# geosync - importe les couches partagées via owncloud dans le geoserver
#*/1 * * * * /usr/local/geosync/bin/sync_and_publish_data.sh && /usr/local/geosync/bin/clean_data.sh >> /var/log/georchestra-ouvert/cron_error.log

*/1 * * * * /usr/local/geosync/bin/sync_and_publish_data.sh 1>>/var/log/georchestra-ouvert/cron.log 2>>/var/log/georchestra-ouvert/cron_error.log && /usr/local/geosync/bin/clean_data.sh 1>>/var/log/georchestra-ouvert/cron.log 2>>/var/log/georchestra-ouvert/cron_error.log

